# kanban-task-tracker-vpn.service
[Unit]
Description=Kanban Task Tracker (Docker Compose, VPN)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/zahar/projects/KanbanTaskTracker
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.vpn.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.vpn.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
